# roles/worker/files/worker.service
[Unit]
Description=Worker Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/application
ExecStart=/usr/bin/java -jar /path/to/your/application/worker.jar
Restart=always

[Install]
WantedBy=multi-user.target
